# backend/webhook.service
[Unit]
Description=ChromaAdapt AI Webhook Listener
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/chromaadapt-ai/backend
Environment="WEBHOOK_SECRET=your-secret-here"
Environment="PROJECT_DIR=/opt/chromaadapt-ai"
Environment="DEPLOY_SCRIPT=/opt/chromaadapt-ai/backend/deploy.sh"
ExecStart=/usr/bin/python /opt/chromaadapt-ai/backend/webhook_server.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/webhook.log
StandardError=append:/var/log/webhook.log

[Install]
WantedBy=multi-user.target
